# 20120204NotificationServer2/notifications-server-assembly/bin/thrift-service
#!/bin/sh
#
### BEGIN INIT INFO
# Provides:          avro-thrift
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO

RETVAL=0
prog="notifications-thrift-server"

usage="Usage: thrift-service (start|stop|install|uninstall)"

# if no args specified, show usage
if [ $# -lt 1 ]; then
  echo $usage
  exit 1
fi

JAVA_HOME_REAL=/usr/lib/jvm/java-6-openjdk
if [ -e "$JAVA_HOME" ]; then
echo JAVA_HOME=$JAVA_HOME
JAVA_HOME_REAL=$JAVA_HOME
elif [ -e "$JAVA_HOME_REAL" ]; then
echo JAVA_HOME_REAL=$JAVA_HOME_REAL
else
echo JAVA_HOME_REAL=[$JAVA_HOME_REAL] not found.
echo JAVA_HOME=[$JAVA_HOME] not found!
exit 1
fi

jsvcfile=jsvc
osbit=$(getconf LONG_BIT)
if [ 32 -eq $osbit ]; then
jsvcfile=jsvc32
fi

RLINK=$(readlink -f $0)

NOTIFJSVC_HOME=$(dirname $RLINK)
FILE_NAME=$(basename $RLINK)

# for multi instances adapt those lines.
TMP_DIR=/var/tmp
PID_FILE=$NOTIFJSVC_HOME/notifications-thrift-service.pid

#echo RLINK=$RLINK
echo NOTIFJSVC_HOME=$NOTIFJSVC_HOME
echo FILE_NAME=$FILE_NAME

CLASSPATH=$NOTIFJSVC_HOME/lib/notifications-server-1.0.0.jar:$NOTIFJSVC_HOME/lib/mongo-java-driver-2.6.3.jar:$NOTIFJSVC_HOME/lib/commons-daemon-1.0.8.jar:$NOTIFJSVC_HOME/lib/activation-1.1.jar:$NOTIFJSVC_HOME/lib/activemq-all-5.5.0.jar:$NOTIFJSVC_HOME/lib/aopalliance-1.0.jar:$NOTIFJSVC_HOME/lib/asm-3.3.1.jar:$NOTIFJSVC_HOME/lib/asm-commons-3.3.1.jar:$NOTIFJSVC_HOME/lib/asm-tree-3.3.1.jar:$NOTIFJSVC_HOME/lib/aspectjrt-1.6.11.jar:$NOTIFJSVC_HOME/lib/aspectjweaver-1.6.11.jar:$NOTIFJSVC_HOME/lib/avalon-framework-4.1.3.jar:$NOTIFJSVC_HOME/lib/libthrift-0.7.0.jar:$NOTIFJSVC_HOME/lib/cglib-nodep-2.2.jar:$NOTIFJSVC_HOME/lib/commons-collections-3.1.jar:$NOTIFJSVC_HOME/lib/commons-dbcp-1.2.1.jar:$NOTIFJSVC_HOME/lib/commons-lang-2.1.jar:$NOTIFJSVC_HOME/lib/commons-logging-1.1.jar:$NOTIFJSVC_HOME/lib/commons-pool-1.2.jar:$NOTIFJSVC_HOME/lib/dom4j-1.6.1.jar:$NOTIFJSVC_HOME/lib/hamcrest-core-1.1.jar:$NOTIFJSVC_HOME/lib/icu4j-2.6.1.jar:$NOTIFJSVC_HOME/lib/jackson-core-asl-1.7.3.jar:$NOTIFJSVC_HOME/lib/jackson-core-lgpl-1.8.3.jar:$NOTIFJSVC_HOME/lib/jackson-mapper-asl-1.7.3.jar:$NOTIFJSVC_HOME/lib/jackson-mapper-lgpl-1.8.3.jar:$NOTIFJSVC_HOME/lib/jars.txt:$NOTIFJSVC_HOME/lib/jaxb-api-2.2.4.jar:$NOTIFJSVC_HOME/lib/jaxb-impl-2.2.4.jar:$NOTIFJSVC_HOME/lib/jdom-1.1.jar:$NOTIFJSVC_HOME/lib/jetty-6.1.26.jar:$NOTIFJSVC_HOME/lib/jetty-util-6.1.26.jar:$NOTIFJSVC_HOME/lib/joda-time-1.6.2.jar:$NOTIFJSVC_HOME/lib/json-20080701.jar:$NOTIFJSVC_HOME/lib/junit-4.10.jar:$NOTIFJSVC_HOME/lib/log4j-1.2.15.jar:$NOTIFJSVC_HOME/lib/logkit-1.0.1.jar:$NOTIFJSVC_HOME/lib/mysql-connector-java-5.1.17.jar:$NOTIFJSVC_HOME/lib/netty-3.2.4.Final.jar:$NOTIFJSVC_HOME/lib/notifications-thrift-api-1.2.0.jar:$NOTIFJSVC_HOME/lib/paranamer-2.3.jar:$NOTIFJSVC_HOME/lib/servlet-api-2.3.jar:$NOTIFJSVC_HOME/lib/servlet-api-2.5-20081211.jar:$NOTIFJSVC_HOME/lib/slf4j-api-1.5.8.jar:$NOTIFJSVC_HOME/lib/slf4j-log4j12-1.5.8.jar:$NOTIFJSVC_HOME/lib/snappy-java-1.0.3.2.jar:$NOTIFJSVC_HOME/lib/spring-aop-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-asm-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-aspects-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-beans-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-context-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-context-support-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-core-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-expression-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-jdbc-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-jms-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-oxm-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-test-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-tx-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/spring-web-3.0.6.RELEASE.jar:$NOTIFJSVC_HOME/lib/stax-api-1.0-2.jar:$NOTIFJSVC_HOME/lib/stax-api-1.0.1.jar:$NOTIFJSVC_HOME/lib/velocity-1.7.jar:$NOTIFJSVC_HOME/lib/wstx-asl-3.2.7.jar:$NOTIFJSVC_HOME/lib/xalan-2.6.0.jar:$NOTIFJSVC_HOME/lib/xercesImpl-2.0.2.jar:$NOTIFJSVC_HOME/lib/xml-apis-1.0.b2.jar:$NOTIFJSVC_HOME/lib/xmlParserAPIs-2.6.2.jar:$NOTIFJSVC_HOME/lib/xom-1.0.jar:$NOTIFJSVC_HOME/lib/xpp3-1.1.4c.jar

case "$1" in
	start)
		if [ -e $PID_FILE ]; then
		  if kill -0 `cat $PID_FILE` > /dev/null 2>&1; then
		    echo thrift-service running as process `cat $PID_FILE`.  Stop it first.
		    exit 1
		  fi
		fi
		#
		# Start Notifications AVRO Serivce
		#
		$NOTIFJSVC_HOME/$jsvcfile -java-home $JAVA_HOME_REAL -Djava.io.tmpdir=$TMP_DIR -wait 10 -pidfile $PID_FILE -outfile $NOTIFJSVC_HOME/thrift-service.log -errfile '&1' -cp $CLASSPATH com.borqs.information.rpc.service.ThriftServiceLauncher 
		#
		# To get a verbose JVM
		#-verbose 
		# To get a debug of jsvc.
		#-debug 
		exit $?
		;;

	stop)
		if [ -e $PID_FILE ]; then
			if kill -0 `cat $PID_FILE` > /dev/null 2>&1; then
				#
				# Stop Notifications AVRO Serivce
				#
				$NOTIFJSVC_HOME/$jsvcfile -stop -pidfile $PID_FILE com.borqs.information.rpc.service.AvroServiceLauncher
			fi
		else
			echo no thrift-service to stop
		fi
		exit $?
		;;
	install)
		if [ -e "/etc/init.d/"$FILE_NAME ]; then
			echo "service already installed."
		else 
			sudo ln -s $NOTIFJSVC_HOME/$FILE_NAME "/etc/init.d/"$FILE_NAME
			sudo update-rc.d $FILE_NAME defaults
                   	echo "service installed successfull."
		fi
		exit $?
		;;
	uninstall)
		if [ -e "/etc/init.d/"$FILE_NAME ]; then
			sudo update-rc.d -f $FILE_NAME remove
			sudo rm "/etc/init.d/"$FILE_NAME
			echo "service uninstalled successfull."
		else 
                   	echo "service not installed yet."
		fi
		exit $?
		;;
	*)
		echo $usage
		exit 1
		;;
esac
